# templates/launchd/com.nyann.sentinel-aggregate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  nyann multi-repo CI sentinel aggregate daemon — launchd user agent (macOS).
  Generated by bin/sentinel-daemon.sh from this template; do not edit the
  installed copy at ~/Library/LaunchAgents/com.nyann.sentinel-aggregate.plist
  by hand — re-run `/nyann:watch --daemon --aggregate` (or sentinel-daemon.sh
  start --aggregate) to regenerate.

  There is ONE aggregate daemon per user (it fans out across every repo in the
  watch-list), so this agent uses a distinct Label from the per-repo agent.

  Placeholders below are substituted at generate time:
    ${AGGREGATE_BIN}  absolute path to bin/sentinel-aggregate.sh
    ${INTERVAL}       base poll interval in seconds
    ${MAX_RUNTIME}    hard cap in seconds (orphan backstop; default 8h)
    ${STATE_DIR}      state/cache dir
    ${NOTIF_DIR}      notification-queue dir
    ${WATCH_LIST}     watch-list JSON the aggregator fans out across
    ${LOG_DIR}        daemon stdout/stderr log dir

  RunAtLoad starts it once on load. KeepAlive is intentionally OFF: the
  daemon self-terminates at --max-runtime as an orphan backstop, and we do
  NOT want launchd to immediately respawn it past that cap. `stop` unloads
  the agent so it stays down until explicitly restarted.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nyann.sentinel-aggregate</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>${AGGREGATE_BIN}</string>
    <string>--daemon-loop</string>
    <string>--interval</string>
    <string>${INTERVAL}</string>
    <string>--max-runtime</string>
    <string>${MAX_RUNTIME}</string>
    <string>--state-dir</string>
    <string>${STATE_DIR}</string>
    <string>--notif-dir</string>
    <string>${NOTIF_DIR}</string>
    <string>--watch-list</string>
    <string>${WATCH_LIST}</string>
    <string>--supervisor</string>
    <string>launchd</string>
  </array>
  <!--
    launchd agents start with a minimal PATH (often just /usr/bin:/bin), so a
    Homebrew-installed `gh` (/opt/homebrew/bin on Apple silicon, /usr/local/bin
    on Intel) would not be found and the daemon would silently no-op. Set an
    explicit PATH that covers both Homebrew prefixes plus the system dirs.
  -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>${LOG_DIR}/sentinel-aggregate.out.log</string>
  <key>StandardErrorPath</key>
  <string>${LOG_DIR}/sentinel-aggregate.err.log</string>
</dict>
</plist>
