# templates/launchd/com.nyann.sentinel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  nyann CI sentinel daemon — launchd user agent (macOS).
  Generated by bin/sentinel-daemon.sh from this template; do not edit the
  installed copy at ~/Library/LaunchAgents/com.nyann.sentinel.plist by hand —
  re-run `/nyann:watch --daemon` (or sentinel-daemon.sh start) to regenerate.

  Placeholders below are substituted at generate time:
    ${SENTINEL_BIN}   absolute path to bin/ci-sentinel.sh
    ${REPO}           <owner>/<repo> being watched
    ${INTERVAL}       poll interval in seconds
    ${MAX_RUNTIME}    hard cap in seconds (orphan backstop; default 8h)
    ${STATE_DIR}      state/cache dir
    ${NOTIF_DIR}      notification-queue dir
    ${LOG_DIR}        daemon stdout/stderr log dir

  RunAtLoad starts it once on load. KeepAlive is intentionally OFF: the
  daemon self-terminates at --max-runtime as an orphan backstop, and we do
  NOT want launchd to immediately respawn it past that cap. `stop` unloads
  the agent so it stays down until explicitly restarted.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nyann.sentinel</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>${SENTINEL_BIN}</string>
    <string>--daemon-loop</string>
    <string>--repo</string>
    <string>${REPO}</string>
    <string>--interval</string>
    <string>${INTERVAL}</string>
    <string>--max-runtime</string>
    <string>${MAX_RUNTIME}</string>
    <string>--state-dir</string>
    <string>${STATE_DIR}</string>
    <string>--notif-dir</string>
    <string>${NOTIF_DIR}</string>
    <string>--supervisor</string>
    <string>launchd</string>
  </array>
  <!--
    launchd agents start with a minimal PATH (often just /usr/bin:/bin), so a
    Homebrew-installed `gh` (/opt/homebrew/bin on Apple silicon, /usr/local/bin
    on Intel) would not be found and the daemon would silently no-op. Set an
    explicit PATH that covers both Homebrew prefixes plus the system dirs.
  -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>${LOG_DIR}/sentinel.out.log</string>
  <key>StandardErrorPath</key>
  <string>${LOG_DIR}/sentinel.err.log</string>
</dict>
</plist>
